# squashfs-root/etc/init.d/xt_qqvoice
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=43

start()
{
	insmod xt_qqvoice
}

stop()
{
	rmmod ipt_salist
}
