# squashfs-root/etc/hotplug.d/gmac/wan/00-internet-led
#!/bin/sh /etc/rc.common

stop() {
	setled off green internet
}
