# squashfs-root/usr/lib/opkg/info/base-files.list
/etc/hotplug.d/gmac/00-normol
/etc/openwrt_version
/etc/resolv.conf
/etc/init.d/done
/lib/firstboot/01_down_wifi
/lib/firstboot/99_10_no_fo_cleanup
/etc/hotplug2-init.rules
/etc/preinit
/lib/firstboot/50_pivot
/sbin/mount_root
/lib/firstboot/30_reset_copy_rom
/etc/wlan/iNIC_ap.dat
/lib/upgrade/platform.sh
/etc/hotplug2-common.rules
/etc/shells
/lib/preinit/70_initramfs_test
/lib/storage.sh
/lib/preinit/50_choose_console
/lib/firstboot/10_no_fo_clear_overlay
/etc/protocols
/etc/uci-defaults/wlan
/sbin/preupgrade.sh
/lib/preinit/20_check_jffs2_ready
/lib/upgrade/hiwifi.sh
/etc/wlan/RT2860.dat
/lib/preinit/10_indicate_preinit
/etc/hotplug.d/gmac/wan/02-dhcp-operation
/lib/wifi/ralink.sh
/lib/firstboot/99_10_with_fo_cleanup
/etc/hotplug.d/button/00-button
/sbin/sdfunc.sh
/lib/functions/net_update_actions.sh
/etc/config/system
/sbin/powerport.sh
/etc/init.d/normal-mode
/etc/hotplug.d/gmac/wan/00-internet-led
/etc/fstab
/lib/preinit/30_device_fs_daemons
/lib/preinit/90_mount_no_jffs2
/etc/sysctl.conf
/sbin/sysupgrade
/lib/preinit/70_pivot_jffs2_root
/lib/preinit/99_10_run_init
/lib/firstboot/20_has_mini_fo
/lib/firstboot/10_reset_has_mini_fo
/lib/firstboot/10_determine_parts
/sbin/led.sh
/lib/firstboot/40_no_fo_copy_ramoverlay
/lib/preinit/50_indicate_regular_preinit
/etc/banner
/sbin/hwf-config
/lib/preinit/10_essential_fs
/etc/shadow
/etc/rc.common
/lib/preinit/10_indicate_failsafe
/lib/preinit/40_init_shm
/init
/etc/init.d/watchdog
/sbin/sdcheck.sh
/sbin/setled
/lib/preinit/40_mount_jffs2
/lib/preinit/20_device_fs_mount
/sbin/switch_test
/lib/preinit/80_mount_root
/etc/rc.local
/lib/functions/uci-defaults.sh
/etc/init.d/led
/lib/functions/network.sh
/etc/init.d/sysctl
/sbin/rled
/lib/firstboot/30_is_rootfs_mounted
/sbin/hotplug-call
/etc/init.d/usb
/lib/platform.sh
/etc/wlan/MT7610E-V10-FEM-1ANT-20130529.bin
/lib/preinit/40_mount_devpts
/lib/preinit/90_init_console
/etc/services
/lib/preinit/99_10_mount_no_mtd
/lib/preinit/10_check_for_mtd
/lib/upgrade/common.sh
/etc/hotplug.d/block/action.d/varlog-sync2disk
/sbin/sdtest.sh
/etc/uci-defaults/board-spec
/etc/sysupgrade.conf
/lib/firstboot/05_firstboot_skip
/etc/hotplug.d/block/10-mount
/etc/init.d/bridge
/lib/upgrade/keep.d/base-files
/lib/preinit/40_run_failsafe_hook
/etc/group
/etc/uci-defaults/migrate-sysctl
/lib/preinit/30_failsafe_wait
/etc/init.d/leds-stats
/etc/functions.sh
/sbin/firstboot
/bin/login.sh
/etc/hosts
/etc/uci-defaults/migrate-shadow
/etc/hotplug.d/usb/10-usb
/etc/wlan/MT7620_AP_2T2R-4L_V15.BIN
/lib/functions.sh
/lib/preinit/60_init_hotplug
/etc/hotplug.d/block/action.d/10-sysupgrade-log
/sbin/speedtest.sh
/etc/init.d/rcS
/etc/uci-defaults/leds
/sbin/diskspace_chk.sh
/rom/note
/etc/inittab
/etc/openwrt_release
/lib/preinit/05_mount_skip
/etc/config/bridge
/sbin/wifi
/sbin/br
/lib/preinit/04_config_vlan
/etc/profile
/etc/uci-defaults/buttons
/etc/wlan/MT7612E_EEPROM.bin
/lib/preinit/99_10_failsafe_login
/etc/TZ
/lib/firstboot/20_no_fo_mount_jffs
/lib/functions/service.sh
/etc/init.d/umount
/lib/ralink.sh
/lib/preinit/90_restore_config
/etc/hotplug.d/wifi/10-normol
/lib/upgrade/keep.d/base-files-essential
/lib/firstboot/30_no_fo_pivot
/etc/passwd
/etc/hotplug.d/ieee1394/10-ieee1394
/lib/preinit/02_default_set_state
/etc/profile.hiwifi
/etc/mtab
/lib/firstboot/40_copy_ramoverlay
/etc/init.d/boot
/etc/hotplug.d/gmac/wan/01-net-update-act
/sbin/storage_format
/etc/uci-defaults/system_defaults
/lib/functions/boot.sh
/bin/ipcalc.sh
/etc/uci-defaults/hiwifi
/etc/diag.sh
/lib/firstboot/20_reset_clear_jffs
/lib/preinit/03_preinit_do_ralink.sh
/lib/preinit/41_merge_overlay_hooks
/etc/uci-defaults/network
