# squashfs-root/usr/lib/opkg/info/cmagent.list
/etc/ca/gd-class2-root.crt
/usr/lib/cmagent/script.lua
/usr/lib/cmagent/diagnose.sh
/etc/cron/360/market_report_installed_apps
/usr/lib/cmagent/cron.lua
/usr/lib/cmagent/p2p.lua
/usr/lib/cmagent/mobileapp.lua
/usr/lib/market/call_func
/etc/hotplug.d/block/restart.d/60-market
/usr/lib/cmagent/status.lua
/etc/init.d/market
/etc/ca/f081611a.0
/usr/lib/cmagent/credential.lua
/usr/lib/cmagent/firmware.lua
/usr/lib/cmagent/conf.lua
/etc/init.d/cmagent
/usr/lib/cmagent/setup.lua
/usr/lib/cmagent/census.lua
/usr/lib/cmagent/reboot.lua
/etc/hotplug.d/block/restart.d/40-cmagent-reload
/usr/lib/cmagent/twx.lua
/lib/upgrade/keep.d/cmagent
/etc/ca/219d9499.0
/usr/lib/cmagent/cmagent
/usr/lib/cmagent/plugin.lua
