# squashfs-root/usr/lib/opkg/info/hwf-openapi.list
/etc/kproxy.d/openapi.hosts
/usr/lib/lua/openapi/apps/appengine.lua
/etc/nginx/vh.client_openapi_hiwifi_com.conf
/usr/lib/lua/openapi/network/wireless.lua
/usr/lib/lua/openapi/system/sys.lua
/etc/init.d/openapi
/usr/lib/cmagent/openapi.lua
/usr/lib/lua/openapi/network/online.lua
/usr/lib/lua/openapi/network/qos.lua
/usr/lib/lua/openapi/callapi.lua
/usr/lib/lua/openapi/network/common.lua
/etc/config/network-conf
/usr/lib/lua/openapi/network/traffic.lua
/usr/lib/lua/openapi/system/apps.lua
/usr/lib/lua/openapi/utils/utils.lua
/usr/lib/lua/openapi/system/rom.lua
/etc/openapi.d/demo/cloud/permissions
/usr/lib/lua/openapi/authz.lua
/usr/lib/lua/openapi/network/twx.lua
/usr/lib/lua/openapi/system/led.lua
/usr/lib/openapi/cgi-bin/openapi-cgi
/usr/lib/lua/openapi/system/os.lua
