# squashfs-root/usr/lib/opkg/info/inet_chk.list
/usr/bin/inet_chk.sh
/usr/bin/inet_chk_switch
/etc/init.d/inet_chk
