# squashfs-root/usr/lib/opkg/info/kmod-hwf-core.list
/lib/modules/3.3.8/hwf_core.ko
/etc/init.d/kcore
