# squashfs-root/usr/lib/opkg/info/kmod-smartqos.list
/lib/functions/hwf_sqos_tc_actions.sh
/lib/modules/3.3.8/hwf_smartqos.ko
/etc/init.d/smartqos
