# squashfs-root/usr/lib/opkg/info/kproxy-tools.list
/etc/kproxy.d/default.hosts
/usr/sbin/kproxy-conf
/usr/lib/lua/kproxy.lua
/usr/sbin/kproxy
/etc/init.d/kproxy
