# squashfs-root/usr/lib/opkg/info/nginx.list
/etc/nginx/logformat.tw-download.conf
/etc/nginx/switch
/etc/nginx/nginx.conf
/etc/nginx/fastcgi_params
/etc/nginx/tw.client_addr.conf
/etc/nginx/uwsgi_params
/etc/nginx/koi-win
/etc/nginx/turbo_aws_502.html
/etc/nginx/vh.tw-download.conf
/etc/nginx/logformat.all.conf
/etc/nginx/vh.tw-vars.conf
/usr/sbin/nginx
/etc/nginx/win-utf
/etc/nginx/turbo_aws_503.html
/etc/nginx/koi-utf
/etc/nginx/arp.lua
/etc/nginx/vh.tw-netchk.conf
/etc/nginx/upstream.all.conf
/etc/hotplug.d/block/restart.d/20-nginx
/etc/init.d/nginx
/etc/nginx/scgi_params
/etc/nginx/fastcgi.conf
/etc/nginx/turbo_aws_504.html
/etc/nginx/mime.types
/etc/nginx/vh.all.conf
/etc/nginx/upstream.tw-download.conf
/etc/nginx/proxy
/etc/nginx/upstream.mcdownloader.conf
/etc/nginx/upstream.error.conf
/etc/nginx/vh.error.conf
/etc/nginx/turbo_aws_500.html
/etc/nginx/vh.4006024680_com.conf
/etc/nginx/mode
