# squashfs-root/usr/lib/opkg/info/pdnsd.list
/usr/sbin/pdnsd
/etc/init.d/pdnsd
/etc/pdnsd.conf
