# squashfs-root/usr/lib/opkg/info/rsyslog.list
/etc/init.d/rsyslog
/usr/lib/rsyslog/imuxsock.so
/usr/lib/rsyslog/lmtcpsrv.so
/usr/lib/rsyslog/lmnetstrms.so
/usr/lib/rsyslog/immark.so
/etc/logrotate.d/rsyslog
/usr/lib/rsyslog/lmregexp.so
/usr/lib/rsyslog/imtcp.so
/usr/lib/rsyslog/omhifwd.so
/usr/sbin/rsyslogd
/usr/lib/rsyslog/imudp.so
/usr/lib/rsyslog/lmtcpclt.so
/usr/lib/rsyslog/lmnet.so
/usr/lib/rsyslog/lmnsd_ptcp.so
/usr/lib/rsyslog/imklog.so
/usr/lib/rsyslog/imfile.so
/etc/rsyslog.conf
/usr/lib/rsyslog/lmzlibw.so
