# squashfs-root/usr/lib/opkg/info/syslog-ng3.list
/etc/init.d/syslog-ng
/etc/syslog-ng.conf_sample
/usr/sbin/syslog-ng
