# squashfs-root/usr/lib/opkg/info/wisp_chk.list
/etc/init.d/wisp_chk
/usr/bin/wisp_chk.sh
