# squashfs-root/usr/lib/opkg/info/xt_qqvoice.list
/usr/lib/iptables/libxt_qqvoice.so
/etc/init.d/xt_qqvoice
