# squashfs-root/etc/crontabs/root
*/10 * */1 * * /usr/sbin/logrotate /etc/logrotate.conf
*/1 * * * * lua /etc/app/traffic_his_2d_script.lua #traffic_his_2d
3 * * * * lua /etc/app/traffic_his_2d_clean.lua #traffic_his_2d
