# squashfs-root/etc/hotplug.d/wifi/device_online
#!/bin/sh

action="$ACTION" 
vap="$VAP"
mac="$MACADDR"
rssi="$RSSI"

if [ $action == 'assoc' ]; then
    lua /etc/app/wifi_event.lua "$action" "$vap" "$mac" "$rssi"
fi
